# etc/suricata/crontab
* * * * * sudo logrotate /etc/logrotate.d/suricata
