# etc/bitbridge6.service
[Unit]
Description=Firewalla BitBridge6
After=network.target

[Service]
User=pi
Restart=on-failure
RestartSec=5
WorkingDirectory=/home/pi/firewalla/bin
ExecStart=/bin/bash /home/pi/firewalla/bin/bitbridge6.sh

[Install]
WantedBy=multi-user.target
