# etc/firewalla.service
[Unit]
Description=firewalla
After=dhcpcd.service redis-server.service
Requires=network-online.target redis-server.service
After=fireupgrade.service
#Wants=brofish.service

[Service]
User=pi
Type=oneshot
ExecStart=/home/pi/firewalla/scripts/main-run
ExecStop=/home/pi/firewalla/scripts/main-stop

[Install]
WantedBy=multi-user.target
