# etc/fireapi.service
[Unit]
Description=Firewalla API
After=network.target

[Service]
User=pi
Restart=always
WorkingDirectory=/home/pi/firewalla/api
ExecStart=/bin/bash -c '/home/pi/firewalla/scripts/start_service.sh FireApi >> /home/pi/.forever/api.log  2>&1'

[Install]
WantedBy=multi-user.target
