# platform/purple/files/fcc.timer
[Unit]
Description=Daily Fire Cron Check

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target