# platform/purple/files/fcc.service
[Unit]
Description=Fire Cron Check

[Service]
Type=oneshot
ExecStart=/home/pi/firewalla/scripts/fcc.sh

